# prototype/boards/qemu-common/rootfs_overlay/etc/udev/rules.d/99-memory-hotplug.rules
SUBSYSTEM=="memory", ACTION=="online", RUN+="/usr/local/bin/pg-retune.sh"
